# victoriatraces/ansible/roles/victoriatraces/templates/victoriatraces.service.j2
[Unit]
Description=VictoriaTraces - Distributed Tracing Storage
Documentation=https://docs.victoriametrics.com/victoriatraces/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ victoriatraces_user }}
Group={{ victoriatraces_group }}

ExecStart={{ victoriatraces_bin_path }}/victoria-traces-prod \
    -storageDataPath={{ victoriatraces_data_path }} \
    -httpListenAddr=:{{ victoriatraces_http_port }} \
    -otlpGRPCListenAddr=:{{ victoriatraces_grpc_port }} \
    -otlpGRPC.tls=false

Restart=always
RestartSec=5

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ victoriatraces_data_path }}

StandardOutput=journal
StandardError=journal
SyslogIdentifier=victoriatraces

[Install]
WantedBy=multi-user.target
