# mihomo/monitoring/templates/mihomo-metrics.service.j2
[Unit]
Description=Mihomo Metrics Collector
Documentation=https://github.com/MetaCubeX/mihomo-meta
After=mihomo.service victoriametrics.service
Requires=mihomo.service
BindsTo=victoriametrics.service

[Service]
Type=oneshot
User=nobody
Group=nogroup

# Environment variables
Environment="MIHOMO_API=http://127.0.0.1:9090"
Environment="MIHOMO_SECRET={{ mihomo_secret }}"
Environment="VICTORIAMETRICS_URL=http://127.0.0.1:8428"
Environment="SERVER_LABEL=gtr"

# ExecStart
ExecStart=/usr/local/bin/mihomo-metrics.sh

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/tmp

# Resource limits
TimeoutStartSec=60
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
