# mihomo/ansible/roles/mihomo/templates/mihomo.service.j2
[Unit]
Description=Mihomo Daemon, Another Clash Kernel
Documentation=https://wiki.metacubex.one
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mihomo
Group=mihomo
LimitNPROC=500
LimitNOFILE=1000000
# Capabilities for TPROXY/TUN mode
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE
Restart=always
RestartSec=5s
ExecStart={{ mihomo_bin_path }}/mihomo -d {{ mihomo_config_path }}
ExecReload=/bin/kill -HUP $MAINPID

# Security hardening
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ mihomo_config_path }} {{ mihomo_home_path }} /var/log

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mihomo

[Install]
WantedBy=multi-user.target
