# edge/ansible/roles/deploy-mutex/templates/deploy-mutex.target.j2
[Unit]
Description=Deployment Mutex Gate
# When busy target is active, new deployments should wait
RefuseManualStart=no
RefuseManualStop=no
